# .config_new/autostart/redshift-gtk.desktop
[Desktop Entry]
X-GNOME-Autostart-enabled=false
Comment[es]=Herramienta para el ajuste de la temperatura de color
StartupNotify=true
Name[de]=Redshift
Comment[de]=Farbtemperatur Anpassungswerkzeug
Name[es]=Redshift
Name[lt]=Redshift
Name[hu]=Redshift
Icon=redshift
Name[tr]=Redshift
Categories=Utility;
GenericName[tr]=Renk sıcaklığı ayarı
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
Comment[hu]=Színhőmérséklet beállító eszköz
GenericName[pl]=Dostosowanie temperatury barwowej
GenericName[hu]=Szín hőmérséklet beállítás
Version=1.0
Comment[tr]=Renk sıcaklığı ayarlama aracı
GenericName[lt]=Spalvos temperatūros reguliavimas
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Comment[da]=Justeringsværktøj for farvetemperatur
GenericName[cs]=Přizpůsobení barevné teploty
GenericName[da]=Justering af farvetemperatur
Comment=Color temperature adjustment tool
Name[ru]=Redshift
Terminal=false
GenericName[de]=Farbtemperatur Anpassung
Hidden=true
Name[fr]=Redshift
Comment[zh_CN]=色温调节工具
Name[pt_BR]=Redshift
GenericName[it]=Regolazione della temperatura del colore
Name[zh_CN]=红移
GenericName[nl]=Kleurtemperatuur aanpassing
GenericName[zh_CN]=色温调节
Name[it]=Redshift
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
Name[pl]=Redshift
Comment[ru]=Инструмент регулирования цветовой температуры
GenericName[es]=Ajuste de la temperatura de color
Exec=redshift-gtk
Comment[fr]=Outil de réglage de la température de couleur
GenericName[fr]=Réglage de la température de couleur
Comment[it]=Strumento per la regolazione della temperatura del colore
GenericName[ru]=Настройка цветовой температуры
Name=Redshift
GenericName=Color temperature adjustment
Type=Application

